# Docker/Base/.config/bioportainer.conf
[program:bioportainer]
command=python3 /workdir/script-server/launcher.py
stderr_logfile=/.config/bioportainer.log
stdout_logfile=/.config/bioportainer.log
autostart=true
