# Docker/Base/.config/bioportainer.conf
[program:bioportainer]
command=/workdir/tools/bioportainer
stderr_logfile=/.config/bioportainer.log
stdout_logfile=/.config/bioportainer.log
autostart=true
