# Docker/Base/.config/guiserver.conf
[program:guiserver]
command=python3 /root/tools/script-server/launcher.py
stderr_logfile=/.config/guiserver.log
stdout_logfile=/.config/guiserver.log
autostart=true
