# Docker/Base/.config/portus.conf
[program:portus]
command=
stderr_logfile=/.config/portus.log
stdout_logfile=/.config/portus.log
autostart=true
