# Docker/Base/.config/portus.conf
[program:portus]
command=/usr/local/bin/portus.sh
stderr_logfile=/.config/portus.log
stdout_logfile=/.config/portus.log
autostart=true
