# Modules/bioportainer-panel/.config/bioportainer.conf
[program:bioportainer]
command=/root/tools/bioportainer
stderr_logfile=/.config/bioportainer.log
stdout_logfile=/.config/bioportainer.log
autostart=true
