# Modules/bioportainer-pipeline-runner/.config/docker.conf
[program:docker]
user=root
autostart=true
command=/usr/local/bin/wrapdocker
redirect_stderr=true
stdout_logfile=/.config/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
