# Modules/bioportainer-pipeline-runner/.config/gui.conf
[program:gui]
user=root
command=python3 /GUI/launcher.py
stderr_logfile=/GUI/guiserver.log
stdout_logfile=/GUI/guiserver.log
autostart=true
