# Modules/bioportainer-pipeline-runner/.config/jupyter.conf
[program:JupyterNotebook]
user=www-runner
command=start-notebook.sh
stderr_logfile=/.config/jupyter_err.log
stdout_logfile=/.config/jupyter.log
